# supervisor.conf
[supervisord]
nodaemon=true

[program:xrdp-sesman]
command=/usr/sbin/xrdp-sesman --nodaemon

[program:xrdp]
command=/usr/sbin/xrdp -nodaemon